# sway/.config/profile.d/80-sway-env.sh
export SWAY_CONFIG_DIR="${XDG_CONFIG_HOME:-$HOME/.config}/sway"
export SWAY_SCRIPTS_DIR="$SWAY_CONFIG_DIR/scripts"
export SWAY_RUNTIME_DIR="${XDG_RUNTIME_DIR:-/var/run/user/$UID}/sway"
